# etc/supervisor-shove.conf
[program:foo]
command=/usr/bin/shove
environment=PYTHONPATH="/usr/lib/python2.7/site-packages/:/usr/lib/python2.7/dist-packages",SHOVE_SETTINGS_FILE=
